# dist/freebsd/mastodon-web
#!/bin/sh

# PROVIDE: mastodon_web
# REQUIRE: LOGIN

# Add the following line to /etc/rc.conf to enable `mastodon_web':
#
#mastodon_enable="YES"

. /etc/rc.subr

name="mastodon_web"
rcvar="mastodon_enable"

load_rc_config "mastodon"
: ${mastodon_enable="NO"}
: ${mastodon_path=""}
: ${mastodon_user=""}

export PATH=${PATH}:/usr/local/bin:/usr/local/sbin
export LANG=ja_JP.UTF-8
export LC_ALL=ja_JP.UTF-8

start_cmd=${name}_start
stop_cmd=${name}_stop

mastodon_web_start() {
  cd $mastodon_path
  sudo -u $mastodon_user zsh -c 'RAILS_ENV=production bundle exec rails server -u puma | logger -t mastodon_web &'
}

mastodon_web_stop() {
  pkill -U $mastodon_user -f puma
}

run_rc_command "$1"
