# dist/freebsd/mastodon-streaming
#!/bin/sh

# PROVIDE: mastodon_streaming
# REQUIRE: LOGIN

# Add the following line to /etc/rc.conf to enable `mastodon_streaming':
#
#mastodon_enable="YES"

. /etc/rc.subr

name="mastodon_streaming"
rcvar="mastodon_enable"

load_rc_config "mastodon"
: ${mastodon_enable="NO"}
: ${mastodon_path="/usr/local/www/mastodon"}
: ${mastodon_user="mastodon"}
: ${mastodon_streaming_bind="0.0.0.0"}
: ${mastodon_env="production"}

export PATH=${PATH}:/usr/local/bin:/usr/local/sbin
export LANG=ja_JP.UTF-8
export LC_ALL=ja_JP.UTF-8
export BIND=$mastodon_streaming_bind
export NODE_ENV=$mastodon_env

start_cmd=${name}_start
stop_cmd=${name}_stop

mastodon_streaming_start() {
  cd $mastodon_path
  sudo -u $mastodon_user sh -c 'npm start | logger -t mastodon_streaming &'
}

mastodon_streaming_stop() {
  pkill -U $mastodon_user -f streaming
}

run_rc_command "$1"
